# os/systemd/abacus-selftest.service
[Unit]
Description=Abacus Self-Test (testinfra)
After=abacus-init.service docker.service
Before=abacus-stack.service
Requires=abacus-init.service docker.service

[Service]
Type=oneshot
ExecStart=/opt/abacus-appliance/selftest/run_selftest.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
